# images/resources/service/sv_apache/run
#!/bin/bash

export TZ="$NAGIOS_TIMEZONE"

exec /usr/sbin/httpd -D NO_DETACH
