# images/nagios-host/resources/etc/service/sv_nsca/run
#!/bin/bash

exec $NAGIOS_HOME/libexec/nsca --daemon -c $NAGIOS_HOME/etc/nsca.cfg
